# zxlhhyccc/koolddns/files/etc/hotplug.d/iface/96-koolddns
#!/bin/sh
case "$ACTION" in
	ifup)
    /etc/init.d/koolddns update
		;;
esac
