# zxlhhyccc/luci-app-zerotier-update2/root/etc/hotplug.d/net/41-zerotier
#!/bin/sh

[ "$ACTION" = "add" ] || exit 0
echo "$DEVICENAME" | grep -Eq "^zt" || exit 0

# 双重检查：服务状态 + UCI 配置
if [ -x "/etc/init.d/zerotier" ] && \
   /etc/init.d/zerotier enabled && \
   [ "$(uci -q get zerotier.global.enabled)" = "1" ]; then
   sleep 5
   zerotier-nat start "$DEVICENAME" > /dev/null 2>&1 &
fi
