# zxlhhyccc/luci-theme-leo/root/etc/uci-defaults/30_luci-theme-leo
#!/bin/sh
uci batch <<-EOF
	set luci.themes.AdvancedTomatoleo=/luci-static/leo
                 set luci.themes.AdvancedTomatoleo1=/luci-static/leo1
	commit luci
EOF
exit 0
